# s6-overlay/services.d/nginx/run
#!/usr/bin/with-contenv bash
# Start Nginx

echo "[INFO] Starting Nginx..."
exec nginx -g "daemon off;"
